# lab1.5-user-services/services/simple-service.service
[Unit]
Description=Simple Test Service
After=network.target

[Service]
Type=simple
User=walter
ExecStart=/home/walter/simple-service.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
